# docker/php/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:frankenphp]
command=docker-php-entrypoint --config /etc/frankenphp/Caddyfile --adapter caddyfile
autostart=true
autorestart=true
redirect_stderr = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:queue-work]
command=php artisan queue:work --queue=backups,default --tries=3 --timeout=3600 --sleep=3 --max-jobs=1000
autostart=true
autorestart=true
redirect_stderr = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
